# SupportCompanion/pkgbuild/payload/Library/LaunchDaemons/com.github.macadmins.SupportCompanion.helper.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.github.macadmins.SupportCompanion.helper</string>
	<key>ProgramArguments</key>
	<array>
		<string>/Library/PrivilegedHelperTools/com.github.macadmins.SupportCompanion.helper</string>
	</array>
	<key>MachServices</key>
	<dict>
		<key>com.github.macadmins.SupportCompanion.helper</key>
		<true/>
	</dict>
	<key>RunAtLoad</key>
	<true/>
	<key>Disabled</key>
	<false/>
</dict>
</plist>
